# install/roles/covid19db/templates/covid19db-update.timer
[Unit]
Description=Update COVID19 Datapool database timer

[Timer]
OnBootSec=15min
OnUnitActiveSec=6h 

[Install]
WantedBy=timers.target
